# .config/systemd/user/clipse-scrub.service
[Unit]
Description=Scrub Clipse clipboard history JSON

[Service]
Type=oneshot
ExecStart=%h/.local/share/bin/scrub_clipboard_history.sh
