# .config/systemd/user/pacman-packages-backup.timer
[Unit]
Description=Hourly snapshot of explicitly installed pacman packages

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
